# conf/systemd.service
[Unit]
Description=Trilium notes server
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
Environment=TRILIUM_DATA_DIR=__DATADIR__
Environment=TRILIUM_PORT=__PORT__
Environment="__YNH_NODE_LOAD_PATH__"
Environment=NODE_ENV=production
WorkingDirectory=__FINALPATH__/
ExecStart=__YNH_NODE__ __FINALPATH__/src/www

[Install]
WantedBy=multi-user.target
